# runasroot/deb-builder/sideloader.service
[Unit]
Description=System Sideloader

[Service]
#Type=oneshot #or simple
ExecStart=/usr/local/hcdd/sideloader

[Install]
WantedBy=multi-user.target
